# rpm/etc/profile.d/nidas.csh
if ( "${path}" !~ */opt/nidas/bin* ) then
    # Uncomment this if you want nidas programs in the path for all login users
    # set path = ( $path /opt/nidas/bin )
endif
