# src/nidas/apps/nidas_udp_relay/eol_relay.crontab
PROJECT=LOTOS

*/10 * * * *    $HOME/isfs/scripts/nidas_eol_relay/check_nidas_udp_relay.sh start
